# misc/memoji-supervisor.conf
[program:memoji-webserver]
command=/usr/bin/node /home/ubuntu/me-moji/memoji-server.js -p 80
directory=/home/ubuntu/me-moji/
numprocs=1
autostart=true