# misc/supervisor.conf
[program:memoji-webserver]
command=/usr/bin/node /home/ubuntu/me-moji/memoji-server.js -p 80
numprocs=1
autostart=true